# systemd_user/mako.service
[Unit]
Description=Mako notification daemon
PartOf=sway-session.target


[Service]
ExecStart=/usr/bin/mako


[Install]
WantedBy=sway-session.target